# scripts/mainnet/systemd/csd.service
[Unit]
Description=Compute Substrate mainnet node
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/compute-substrate
Environment=ROOT_DIR=/home/ubuntu/compute-substrate
Environment=BIN_PATH=/home/ubuntu/compute-substrate/target/release/csd
Environment=GENESIS_PATH=/home/ubuntu/compute-substrate/genesis.bin
Environment=DATADIR=/home/ubuntu/.csd/mainnet-node
Environment=RPC_ADDR=0.0.0.0:8789
Environment=P2P_LISTEN=/ip4/0.0.0.0/tcp/17999
Environment=RUST_LOG=info
# Set BOOTNODES here before enabling on real mainnet:
Environment=BOOTNODES=
ExecStart=/home/ubuntu/compute-substrate/scripts/mainnet/run-node.sh
Restart=always
RestartSec=5
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
